# dome.target
[Unit]
Description=Dome
Requires=default.target
After=default.target
