# service/domed.service
[Unit]
Description=Dome daemon

[Service]
Type=dbus
BusName=ru.dome
User=root
Group=root
ExecStart=/usr/bin/domed
