# core/dome_core.service
[Unit]
Description=Dome core

[Service]
Type=simple
WorkingDirectory=%h/dome
ExecStart=/usr/bin/dome_core -c /home/pi/dome/core.json -l d
